# tg-shorts-bot.service
[Unit]
Description=Telegram Shorts Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/tg_shorts_bot
Environment="PATH=/home/YOUR_USERNAME/tg_shorts_bot/venv/bin"
ExecStart=/home/YOUR_USERNAME/tg_shorts_bot/venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
